# .config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Name=Polybar
Exec=polybar
X-GNOME-Autostart-Delay=2
Comment[en_US]=StartupPolybar
